# littleant.service
[Unit]
Description=LittleAnt V14 - AI Intent Execution System
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/littleant
ExecStart=/usr/bin/python3 /littleant/run.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
